# gnu/.config/systemd/user/i3-session.target
[Unit]
Description=i3 graphical session
BindsTo=graphical-session.target
