# gnu/.config/systemd/user/new-screenshots.path
[Unit]
Description=Detect new screenshots
PartOf=graphical-session.target

[Path]
PathChanged=/home/kpg/Zrzuty/
MakeDirectory=true
DirectoryMode=700
