# gnu/.config/systemd/user/new-screenshots.path
[Unit]
Description=Detect new screenshots
PartOf=graphical-session.target

[Path]
PathChanged=%h/Zrzuty/
MakeDirectory=true
DirectoryMode=700
TriggerLimitBurst=0
