# gnu/.config/systemd/user/pasystray.service
[Unit]
Description=Pulseaudio tray applet

[Service]
ExecStart=/usr/bin/env pasystray
Restart=always
