# gnu/.config/systemd/user/xsettingsd-conf.path
[Unit]
Description=Reapply xsettingsd.conf
PartOf=graphical-session.target

[Path]
PathChanged=%h/.config/xsettingsd/xsettingsd.conf
